# cache/config/supervisor.conf
[supervisord]
nodaemon = true

[program:memcached]
command = /usr/local/bin/memcached logfile /var/log/memcached/memcached.log
user = memcache
autostart = true
