# database/config/supervisor.conf
[supervisord]
nodaemon = true

[program:percona]
command = /docker-entrypoint.sh mysqld
user = root
autostart = true
